# deploy/systemd/onlinetest-deploy-hook.service
[Unit]
Description=OnlineTest deploy webhook (POST, localhost only)
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/bin/bash /var/www/onlinetest/deploy/deploy-hook-runner.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
